# roles/nginx-lets-encrypt/files/acmetool/acmetool.cron
SHELL=/bin/sh
PATH=/bin:/sbin:/usr/bin:/usr/sbin:/usr/local/bin:/usr/local/sbin
MAILTO=root
6 11 * * * root /usr/local/bin/acmetool --batch reconcile
